# root/etc/supervisord.d/wechat-monitor.conf
[program:wechat-monitor]
command=wechat-monitor
user=app
autorestart=true
stdout_logfile=/tmp/wechat-monitor.log
stderr_logfile=/tmp/wechat-monitor.log