# wechat-forward/etc/supervisord.d/wechat-forward.conf
[program:wechat-forward]
command=wechat-forward
user=app
autorestart=true
stdout_logfile=/tmp/wechat-forward.log
stderr_logfile=/tmp/wechat-forward.log
